# tinkerforge-brickd/supervisord.conf
[supervisord]
nodaemon=true

[program:brickd]
command=/etc/init.d/brickd start
